# files/system/usr/lib/systemd/system/secure-ai-enable-diffusion.path
[Unit]
Description=Watch for diffusion runtime install request

[Path]
PathExists=/run/secure-ai-ui/diffusion-request

[Install]
WantedBy=multi-user.target
