# files/system/usr/lib/systemd/system/secure-ai-dns-leak-check.timer
[Unit]
Description=Periodic DNS Leak Detection (every 60 minutes)

[Timer]
OnBootSec=2min
OnUnitActiveSec=60min
RandomizedDelaySec=120
Persistent=true

[Install]
WantedBy=timers.target
