# files/system/usr/lib/systemd/system/secure-ai-health-check.service
[Unit]
Description=Secure AI Post-Update Health Check
After=secure-ai-firstboot.service multi-user.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/etc/greenboot/check/required.d/01-secure-ai-health.sh
TimeoutStartSec=300

DynamicUser=yes
ProtectHome=yes
PrivateTmp=yes

# Kernel protection
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectHostname=yes

# Capability bounding — non-root
CapabilityBoundingSet=

# Syscall filtering
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @mount @clock @debug @swap @reboot @module @cpu-emulation @obsolete
SystemCallArchitectures=native
SystemCallErrorNumber=EPERM

LimitCORE=0

ReadWritePaths=/var/lib/secure-ai/logs /run/secure-ai

[Install]
WantedBy=multi-user.target
