# files/system/usr/lib/systemd/system/secure-ai-search-mediator.service
[Unit]
Description=Secure AI Appliance - Search Mediator (Sanitized Web Search)
After=secure-ai-searxng.service
Wants=secure-ai-searxng.service

[Service]
Type=simple
ExecStart=/usr/libexec/secure-ai/search-mediator
Restart=on-failure
RestartSec=5

Environment=BIND_ADDR=127.0.0.1:8485
Environment=SEARXNG_URL=http://127.0.0.1:8888
Environment=APPLIANCE_CONFIG=/etc/secure-ai/config/appliance.yaml
Environment=POLICY_PATH=/etc/secure-ai/policy/policy.yaml
Environment=AUDIT_DIR=/var/lib/secure-ai/logs

# Filesystem isolation — only talks to SearXNG on localhost
DynamicUser=yes
ReadWritePaths=/var/lib/secure-ai/logs
ReadOnlyPaths=/etc/secure-ai
PrivateTmp=yes
PrivateDevices=yes
ProtectSystem=strict
ProtectHome=yes
ProtectProc=invisible
ProcSubset=pid

# Kernel protection
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectHostname=yes

# Privilege restriction
NoNewPrivileges=yes
RestrictSUIDSGID=yes
MemoryDenyWriteExecute=yes
LockPersonality=yes
RestrictRealtime=yes

# Capability bounding
CapabilityBoundingSet=
AmbientCapabilities=

# Namespace restrictions
RestrictNamespaces=yes
RestrictAddressFamilies=AF_INET AF_UNIX

# Clipboard isolation (M21) — no access to user clipboard
PrivateUsers=yes

# Syscall filtering
SystemCallFilter=@system-service @network-io
SystemCallFilter=~@privileged @resources @mount @clock @debug @swap @reboot @raw-io @module @cpu-emulation @obsolete
SystemCallArchitectures=native
SystemCallErrorNumber=EPERM

# Device cgroup — deny all device access (M43/M5)
DevicePolicy=closed
DeviceAllow=

# Resource limits (M43/M5 tightened)
MemoryMax=256M
CPUQuota=25%
TasksMax=64
LimitNPROC=32
LimitNOFILE=512
LimitCORE=0

[Install]
WantedBy=multi-user.target
