# files/system/usr/lib/systemd/system/secure-ai-integrity-monitor.service
[Unit]
Description=Secure AI Continuous Integrity Monitor (baseline-verified file watcher)
# Starts after attestor and policy engine; monitors all service binaries and policies
After=secure-ai-runtime-attestor.service secure-ai-policy-engine.service network.target

[Service]
Type=simple
ExecStart=/usr/libexec/secure-ai/integrity-monitor
Environment=BIND_ADDR=127.0.0.1:8510
Environment=MONITOR_POLICY_PATH=/etc/secure-ai/policy/integrity-monitor.yaml
Environment=AUDIT_LOG_PATH=/var/lib/secure-ai/logs/integrity-monitor-audit.jsonl
Environment=SERVICE_TOKEN_PATH=/run/secure-ai/service-token

# Filesystem isolation
DynamicUser=yes
LogsDirectory=secure-ai
ReadOnlyPaths=/etc/secure-ai
ReadOnlyPaths=/run/secure-ai
ReadOnlyPaths=/usr/libexec/secure-ai
ReadOnlyPaths=/usr/lib/systemd/system
ReadOnlyPaths=/var/lib/secure-ai/registry
ReadWritePaths=/var/lib/secure-ai/logs
PrivateTmp=yes
PrivateDevices=yes
ProtectSystem=strict
ProtectHome=yes
ProtectProc=invisible
ProcSubset=pid

# Kernel protection
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectHostname=yes

# Privilege restriction
NoNewPrivileges=yes
RestrictSUIDSGID=yes
MemoryDenyWriteExecute=yes
LockPersonality=yes
RestrictRealtime=yes

# Capability bounding
CapabilityBoundingSet=
AmbientCapabilities=

# Namespace restrictions
RestrictNamespaces=yes
RestrictAddressFamilies=AF_UNIX AF_INET

# Private users
PrivateUsers=yes

# PrivateNetwork=no — needs localhost HTTP for API and to read service binaries
PrivateNetwork=no

# Syscall filtering
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources @mount @clock @debug @swap @reboot @raw-io @module @cpu-emulation @obsolete
SystemCallArchitectures=native
SystemCallErrorNumber=EPERM

# Resource limits
MemoryMax=256M
CPUQuota=15%
TasksMax=32
LimitNOFILE=1024
LimitCORE=0

# Production timeouts and restart limits
TimeoutStartSec=30
TimeoutStopSec=15
StartLimitInterval=300
StartLimitBurst=5

Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
