# files/system/usr/lib/systemd/system/secure-ai-diffusion.service
[Unit]
Description=Secure AI Appliance - Diffusion Worker (Image/Video Generation)
After=secure-ai-registry.service
Wants=secure-ai-registry.service

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/secure-ai/services/diffusion-worker/app.py
Restart=on-failure
RestartSec=5

Environment=BIND_ADDR=127.0.0.1:8455
Environment=REGISTRY_DIR=/var/lib/secure-ai/registry
Environment=OUTPUTS_DIR=/var/lib/secure-ai/vault/outputs
Environment=APPLIANCE_CONFIG=/etc/secure-ai/config/appliance.yaml
Environment=MAX_RESOLUTION=2048
Environment=MAX_STEPS=100
EnvironmentFile=-/var/lib/secure-ai/inference.env

# Filesystem isolation — needs GPU access, no network
DynamicUser=yes
ReadWritePaths=/var/lib/secure-ai/vault/outputs
ReadOnlyPaths=/var/lib/secure-ai/registry /etc/secure-ai
PrivateNetwork=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ProtectProc=invisible
ProcSubset=pid

# Kernel protection
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectHostname=yes

# Privilege restriction
NoNewPrivileges=yes
RestrictSUIDSGID=yes
MemoryDenyWriteExecute=no
LockPersonality=yes
RestrictRealtime=yes

# Capability bounding
CapabilityBoundingSet=
AmbientCapabilities=

# Namespace restrictions
RestrictNamespaces=yes
RestrictAddressFamilies=AF_UNIX AF_INET

# Clipboard isolation (M21) — no access to user clipboard
PrivateUsers=yes

# Syscall filtering — needs GPU ioctls for CUDA/ROCm, image write
SystemCallFilter=@system-service @io-event
SystemCallFilter=~@privileged @mount @clock @debug @swap @reboot @module @cpu-emulation @obsolete
SystemCallArchitectures=native
SystemCallErrorNumber=EPERM

# GPU access — NVIDIA (CUDA), AMD (ROCm via /dev/kfd + /dev/dri), Intel (DRI)
DeviceAllow=/dev/nvidia* rw
DeviceAllow=/dev/nvidiactl rw
DeviceAllow=/dev/nvidia-uvm rw
DeviceAllow=/dev/nvidia-uvm-tools rw
DeviceAllow=/dev/dri/* rw
DeviceAllow=/dev/kfd rw
SupplementaryGroups=video render

# Resource limits — generous for image/video generation
MemoryMax=32G
CPUQuota=200%
TasksMax=256
LimitNOFILE=4096
LimitCORE=0

[Install]
WantedBy=multi-user.target
