# files/system/usr/lib/systemd/system/secure-ai-boot-verify.service
[Unit]
Description=Secure AI Boot Chain Integrity Verification
After=local-fs.target
Before=secure-ai-registry.service secure-ai-ui.service
DefaultDependencies=no

[Service]
Type=oneshot
ExecStart=/usr/libexec/secure-ai/verify-boot-chain.sh
RemainAfterExit=yes

# Needs access to /sys/firmware/efi, /dev/tpmrm0, /boot
ProtectSystem=full
ProtectHome=yes
PrivateTmp=yes
PrivateNetwork=yes

# Kernel protection
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectHostname=yes

# Capability bounding — root oneshot, no special caps needed
CapabilityBoundingSet=

# Syscall filtering
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @mount @clock @debug @swap @reboot @module @cpu-emulation @obsolete
SystemCallArchitectures=native
SystemCallErrorNumber=EPERM

# Write results to logs
ReadWritePaths=/var/lib/secure-ai/logs /run/secure-ai

[Install]
WantedBy=multi-user.target
