# files/system/usr/lib/systemd/system/secure-ai-searxng.service
[Unit]
Description=Secure AI Appliance - SearXNG Metasearch Engine (Tor-Routed)
After=secure-ai-tor.service
Wants=secure-ai-tor.service

[Service]
Type=simple
Environment=SEARXNG_SETTINGS_PATH=/etc/secure-ai/searxng/settings.yml
ExecStart=/usr/bin/python3 -m searx.webapp
Restart=on-failure
RestartSec=5

# Sandboxing — needs network only to reach Tor SOCKS5 on localhost
DynamicUser=yes
ProtectSystem=strict
ReadOnlyPaths=/etc/secure-ai/searxng
ProtectHome=yes
PrivateTmp=yes
NoNewPrivileges=yes
RestrictSUIDSGID=yes

# Kernel protection
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectHostname=yes

# Capability bounding — non-root
CapabilityBoundingSet=

# SearXNG connects to Tor on localhost, which connects to the internet
RestrictAddressFamilies=AF_INET AF_UNIX

# Clipboard isolation (M21) — no access to user clipboard
PrivateUsers=yes

# No GPU, no devices
PrivateDevices=yes

# Syscall filtering
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @mount @clock @debug @swap @reboot @module @cpu-emulation @obsolete
SystemCallArchitectures=native
SystemCallErrorNumber=EPERM

# Resource limits
MemoryMax=512M
TasksMax=128

# Production timeouts
TimeoutStopSec=15

[Install]
WantedBy=multi-user.target
