# files/system/usr/lib/systemd/system/secure-ai-airlock.service
[Unit]
Description=Secure AI Online Airlock (sanitized egress proxy)
After=secure-ai-tool-firewall.service
Requires=secure-ai-tool-firewall.service

[Service]
Type=simple
ExecStart=/usr/libexec/secure-ai/airlock
Environment=BIND_ADDR=127.0.0.1:8490
Environment=POLICY_PATH=/etc/secure-ai/policy/policy.yaml
Environment=SOURCES_ALLOWLIST_PATH=/etc/secure-ai/policy/sources.allowlist.yaml
Environment=AUDIT_LOG_PATH=/var/lib/secure-ai/logs/airlock-audit.jsonl

# Filesystem isolation
DynamicUser=yes
LogsDirectory=secure-ai
ReadOnlyPaths=/etc/secure-ai
ReadWritePaths=/var/lib/secure-ai/logs
# NOTE: NO PrivateNetwork — airlock is the ONLY service with outbound access
PrivateTmp=yes
PrivateDevices=yes
ProtectSystem=strict
ProtectHome=yes
ProtectProc=invisible
ProcSubset=pid

# Kernel protection
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectHostname=yes

# Privilege restriction
NoNewPrivileges=yes
RestrictSUIDSGID=yes
MemoryDenyWriteExecute=yes
LockPersonality=yes
RestrictRealtime=yes

# Capability bounding — no special caps needed
CapabilityBoundingSet=
AmbientCapabilities=

# Namespace restrictions
RestrictNamespaces=yes
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6

# Clipboard isolation (M21) — no access to user clipboard
PrivateUsers=yes

# Syscall filtering
SystemCallFilter=@system-service @network-io
SystemCallFilter=~@privileged @resources @mount @clock @debug @swap @reboot @raw-io @module @cpu-emulation @obsolete
SystemCallArchitectures=native
SystemCallErrorNumber=EPERM

# Resource limits — modest, only proxies small requests
MemoryMax=256M
CPUQuota=25%
TasksMax=32
LimitNOFILE=512
LimitCORE=0

# Disabled by default — enable only when online augmentation is needed.
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
