# files/system/usr/lib/systemd/system/secure-ai-registry.service
[Unit]
Description=Secure AI Trusted Registry (artifact manifest + model store)
After=local-fs.target secure-ai-firstboot.service
Wants=secure-ai-firstboot.service

[Service]
Type=simple
ExecStart=/usr/libexec/secure-ai/registry
Environment=BIND_ADDR=127.0.0.1:8470
Environment=REGISTRY_DIR=/var/lib/secure-ai/registry
Environment=REGISTRY_LOCK_PATH=/etc/secure-ai/policy/models.lock.yaml
Environment=AUDIT_LOG_PATH=/var/lib/secure-ai/logs/registry-audit.jsonl
Environment=SERVICE_TOKEN_PATH=/run/secure-ai/service-token

# Filesystem isolation
DynamicUser=yes
StateDirectory=secure-ai/registry
LogsDirectory=secure-ai
ReadOnlyPaths=/etc/secure-ai
ReadOnlyPaths=/run/secure-ai
ReadWritePaths=/var/lib/secure-ai/registry
ReadWritePaths=/var/lib/secure-ai/logs
PrivateNetwork=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectSystem=strict
ProtectHome=yes
ProtectProc=invisible
ProcSubset=pid

# Kernel protection
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectHostname=yes

# Privilege restriction
NoNewPrivileges=yes
RestrictSUIDSGID=yes
MemoryDenyWriteExecute=yes
LockPersonality=yes
RestrictRealtime=yes

# Capability bounding
CapabilityBoundingSet=
AmbientCapabilities=

# Namespace restrictions
RestrictNamespaces=yes
RestrictAddressFamilies=AF_UNIX

# Clipboard isolation (M21) — no access to user clipboard
PrivateUsers=yes

# Syscall filtering
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources @mount @clock @debug @swap @reboot @raw-io @module @cpu-emulation @obsolete
SystemCallArchitectures=native
SystemCallErrorNumber=EPERM

# Resource limits
MemoryMax=512M
CPUQuota=50%
TasksMax=64
LimitNOFILE=1024
LimitCORE=0

# Production timeouts and restart limits
TimeoutStartSec=30
TimeoutStopSec=15
StartLimitInterval=300
StartLimitBurst=5

Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
