# files/system/usr/lib/systemd/system/secure-ai-audit-verify.timer
[Unit]
Description=Secure AI Audit Log Chain Verification Timer (default: every 30 minutes)

[Timer]
OnBootSec=5min
OnUnitActiveSec=30min
RandomizedDelaySec=60
Persistent=true

[Install]
WantedBy=timers.target
