# files/system/usr/lib/systemd/system/secure-ai-tor.service
[Unit]
Description=Secure AI Appliance - Tor Anonymous Proxy
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStartPre=/bin/mkdir -p /var/lib/secure-ai/tor /var/lib/secure-ai/logs
ExecStart=/usr/bin/tor -f /etc/secure-ai/tor/torrc
Restart=on-failure
RestartSec=10

# Sandboxing — Tor needs network but nothing else
DynamicUser=yes
StateDirectory=secure-ai/tor
LogsDirectory=secure-ai
ProtectSystem=strict
ReadWritePaths=/var/lib/secure-ai/tor /var/lib/secure-ai/logs
ReadOnlyPaths=/etc/secure-ai/tor
ProtectHome=yes
PrivateTmp=yes
NoNewPrivileges=yes
RestrictSUIDSGID=yes

# Kernel protection
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectHostname=yes

# Tor needs network (it IS the network proxy)
# But restrict to IPv4/IPv6 TCP only
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX

# Clipboard isolation (M21) — no access to user clipboard
PrivateUsers=yes

# No GPU, no devices
PrivateDevices=yes

# Syscall filtering
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @mount @clock @debug @swap @reboot @module @cpu-emulation @obsolete
SystemCallArchitectures=native
SystemCallErrorNumber=EPERM

# Resource limits
MemoryMax=256M
TasksMax=64

# Production timeouts
TimeoutStopSec=15

[Install]
WantedBy=multi-user.target
