# files/system/usr/lib/systemd/system/secure-ai-dns-leak-check.service
[Unit]
Description=Secure AI DNS Leak Detection Check
After=network-online.target secure-ai-tor.service
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/libexec/secure-ai/check-dns-leak.sh

# Sandboxing
DynamicUser=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
NoNewPrivileges=yes

# Kernel protection
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectHostname=yes

# Capability bounding — non-root
CapabilityBoundingSet=

LimitCORE=0

# Needs network to test DNS + Tor connectivity
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX

# Write results to logs
ReadWritePaths=/var/lib/secure-ai/logs
ReadOnlyPaths=/etc/secure-ai

# Syscall filtering
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @mount @clock @debug @swap @reboot @module @cpu-emulation @obsolete
SystemCallArchitectures=native
SystemCallErrorNumber=EPERM

# Resource limits
MemoryMax=64M
TasksMax=16
