# files/system/usr/lib/systemd/system/secure-ai-clipboard-clear.timer
[Unit]
Description=Periodic Clipboard Clear (every 60 seconds)
# Prevents sensitive data from persisting on the clipboard between uses

[Timer]
OnBootSec=2min
OnUnitActiveSec=60s
RandomizedDelaySec=5

[Install]
WantedBy=timers.target
