# files/system/usr/lib/systemd/system/secure-ai-apply-profile.path
[Unit]
Description=Watch for profile change request from UI

[Path]
PathExists=/run/secure-ai-ui/profile-request

[Install]
WantedBy=multi-user.target
